# docker/sshd/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -e
autorestart=true
stdout_logfile=/var/log/sshd/sshd.stdout
stderr_logfile=/var/log/sshd/sshd.stderr
